# crons/cron.daily/fix_failed_syncing
#!/bin/bash
cd /opt/apphomes/hubplus/pinax-env/hubplus/
source ../bin/activate

python manage.py fix_failed_syncing > /dev/null
