# crons/cron.hourly/reindex
#!/bin/bash
cd /opt/apphomes/plusdev/hubplus/
source ../bin/activate
python manage.py update_index > /dev/null
